# system_files/usr/lib/systemd/system/systemd-suspend.service.d/10-armada-fake-suspend.conf
[Unit]
Description=Armada Fake Suspend

[Service]
# Real suspend hangs; run fake-suspend, which blocks until woken.
ExecStart=
ExecStart=/usr/libexec/armada/fake-suspend sleep
TimeoutStartSec=infinity
